# release/src-rt-6.x.4708/router/apcupsd/platforms/debian/packageinfo/postrm
#!/bin/sh -e

if [ "$1" = "purge" ] ; then
    update-rc.d apcupsd remove >/dev/null
    rm -f /etc/apcupsd/apcupsd.conf
    rm -f /etc/apcupsd/apcupsd.status
    rm -f /etc/apcupsd/apcupsd.events
    rm -f /etc/apcupsd
fi
